# imgserv.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='imgserv'>
  <service name='site/imgserv' type='service' version='1'>
    <create_default_instance enabled='true' />
    <single_instance />
    <dependency name='everything' grouping='require_all' restart_on='refresh' type='service'>
      <service_fmri value='svc:/milestone/multi-user-server' />
    </dependency>
    <exec_method type='method' name='start' exec='/opt/local/bin/python /opt/20ft/imgserv/imgserv.py' timeout_seconds='0'>
        <method_context working_directory='/opt/20ft/imgserv' />
    </exec_method>
    <exec_method type='method' name='stop' exec=':kill' timeout_seconds='60' />
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='child'/>
    </property_group>
    <template>
      <common_name>
        <loctext xml:lang='C'>imgadm server</loctext>
      </common_name>
      <description>
        <loctext xml:lang='C'>Lightweight imgadm compatible image server</loctext>
      </description>
    </template>
  </service>
</service_bundle>
